# systemd/pistation-kiosk.service
[Unit]
Description=PiStation Kiosk Dashboard
After=graphical.target pistation-server.service
Wants=pistation-server.service

[Service]
Type=simple
User=pi
Group=pi
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pi/.Xauthority
# Kurze Wartezeit damit X11 sicher bereit ist
ExecStartPre=/bin/bash -c 'for i in $(seq 1 30); do xdpyinfo -display :0 >/dev/null 2>&1 && exit 0; sleep 1; done; echo "FEHLER: X11 Display :0 nicht erreichbar nach 30s" >&2; exit 1'
ExecStart=/home/pi/pi-media-station/services/dashboard-start.sh
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical.target
