# systemd/pistation-player.service
[Unit]
Description=PiStation MPV Player Daemon
# HINWEIS: Dieser Service wird normalerweise NICHT benötigt.
# media-server.py startet mpv über player-controller.py::ensure_running() selbst.
# Nur aktivieren wenn mpv ohne media-server.py laufen soll.
After=network.target

[Service]
Type=simple
User=pi
EnvironmentFile=/home/pi/pi-media-station/.env
Environment=XDG_RUNTIME_DIR=/run/user/1000
ExecStart=/usr/bin/mpv --idle=yes --input-ipc-server=${SOCKET_PATH:-/tmp/mpv-socket} \
          --fullscreen=no --no-terminal --loop-file=no \
          --volume=${MPV_VOLUME:-100} --keep-open=yes
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
# Kein WantedBy — muss manuell aktiviert werden: systemctl enable pistation-player
