# systemd/videoanalyse-wlan-failsafe.service
[Unit]
Description=VideoAnalyse one-time AP manager with 15s rollback
After=network.target dhcpcd.service
Wants=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/wlan-failsafe-manager.sh boot
TimeoutStartSec=45
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
